# deploy/systemd/coordinate-zero-goaccess.timer
[Unit]
Description=Refresh the Coordinate Zero GoAccess report every five minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=coordinate-zero-goaccess.service

[Install]
WantedBy=timers.target
