# deploy/ops/coordinate-zero-ops.cron.example
# /etc/cron.d/coordinate-zero-ops
# Cron uses the server timezone. Set the ECS host to Asia/Shanghai first.
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

*/5 * * * * root /var/www/coordinate-zero/deploy/ops/collect-status.sh >> /var/log/coordinate-zero-ops.log 2>&1
30 3 * * * root /var/www/coordinate-zero/deploy/ops/backup-mysql.sh >> /var/log/coordinate-zero-backup.log 2>&1
